# contrib/systemd/dns-bh_export-file.service
[Unit]
Description=DNS Blacholing Service (export file)

[Service]
Type=oneshot
ExecStart=/opt/dns-bh/bin/export-file -cfg-dir=/opt/dns-bh/etc
